# scripts_all/rc.local
for iface in $(ip -o link | cut -d: -f2 | tr -d ' ' | grep ^eth)
do
   test -f /etc/sysconfig/network-scripts/ifcfg-$iface
   if [ $? -ne 0 ]
   then
       touch /etc/sysconfig/network-scripts/ifcfg-$iface
       echo -e "DEVICE=$iface\nBOOTPROTO=dhcp\nONBOOT=yes" > /etc/sysconfig/network-scripts/ifcfg-$iface
       ifup $iface
   fi
done
exit 0
